# GDS_start_script
#!/bin/bash

### BEGIN INIT INFO
# Provides:          gds
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs $network
# Should-Start:      postgresql
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: start GDS
# Description:  GDS is a deamon who control the git-deployer script
# 		to auto deploy projects
# Author: Guillaume Seigneuret
# CopyRight OmegaCube 2012
### END INIT INFO


GDS_HOME=/home/Git-Deployer
GDS_BIN=$GDS_HOME/GDS.pl
GDS_PID=/var/run/gds.pid
PORT=32337
PATTERN_RECOG="$PORT.*LISTEN.*perl"

# GDS options
#  -d, 		Daemon mode
#  -l log_file, Specify another log file than the one in conf file
#  -p pid_file, Specify another pid file than the one in conf file

[ ! -e "${GDS_BIN}" ] && echo "GDS ERROR - ${GDS_BIN} does not exist" && exit 1

GDS_COMMAND="$GDS_BIN -d -p $GDS_PID"
[ -e $GDS_PID ] && PID=`cat $GDS_PID`


case "$1" in
	start)
		echo -n "Starting Git Deployment Server ... "
		$GDS_COMMAND && echo "[STARTED]" && exit 0
		echo "[FAILED]" && exit 1
		;;
	stop)
		echo -n "Stopping Git Deployment Server ... "
		[ -z $PID ] && echo "Could not retreive GDS PID :(" && exit 1
		kill $PID
		if ps -p $PID > /dev/null
		then
			echo "[FAILED]" && exit 1
		else
			rm -f $GDS_PID
			echo "[STOPPED]"
			exit 0
		fi
		;;
	restart)
		echo -n "Restarting Git Deployment Server "
		$0 stop
		i=0
		exitcode=1
                while [ $exitcode -eq 0 ]
                do
                        result=`netstat -anpe  | grep -E $PATTERN_RECOG`
                        exitcode=$?
			i=`expr $i + 1`	
			[ $i -ge 100 ] && echo "Socket doesn't want to die !" && exit 1
			[ $(($i % 2)) -eq 0 ] && echo -n "."
		done
		echo ""
		$0 start
		exit 0
		;;
	status)
		echo -n "GDS process status ... "
		result=`pgrep GDS`
		GDS_process_status=$?
		netstat -anpe  | grep -e $PATTERN_RECOG > /dev/null 2>&1
		GDS_socket_status=$?

		[ $GDS_process_status -eq 0 ] && echo "Process alive (PID: $result)"
		[ $GDS_process_status -ne 0 ] && echo "No GDS process found :("
		
		echo -n "GDS socket status ... "
		[ $GDS_socket_status -eq 0 ] && echo "Socket alive" && exit 0
		[ $GDS_socket_status -ne 0 ] && echo "Socket not listening !!!" && exit 1
		;;
	*)
		echo "GDS usage: $0 {start|stop|restart|status}"
		exit 1
		;;
esac
